# roles/transit/files/routing.service
[Unit]
Description=Ping IP and start service if successful

[Service]
User=root
ExecStart=/etc/realm/routing.sh